# images/wireless/supervisor_wifi.conf
[program:setup_wifi]
command=python /setup_wifi.py
autorestart=true
