# images/inet/supervisor_inet.conf
[program:startup]
command=python /start_inet.py
autorestart=false
