# images/elk/supervisor_elk.conf
[program:elasticsearch]
command=/bin/bash -c 'service elasticsearch start && sleep 5'
startsec=0
autorestart=false

[program:kibana]
command=/bin/bash -c 'exec /kibana/bin/kibana'

