# images/switch/supervisor_sw.conf
[program:start_sw]
command=python /start_sw.py

[program:lldpd]
command=/bin/bash -c "lldpd -d -ccc -ff -ss -ee -l"
autorestart=true
