# images/victims/supervisor_vics.conf
[program:winbind]
command=/etc/init.d/winbind start
startsec=0
autorestart=false

[program:startup]
command=/bin/bash -c 'sleep 30 && python start_vic.py'
autorestart=true
