# images/samba/supervisor_samba.conf
[program:samba]
command=service samba start
startsecs=0
autorestart=false
